# debian/glance.glance-registry.upstart
description "Glance registry server"
author "Soren Hansen <soren@linux2go.dk>"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

respawn

exec su -c "cd /var/lib/glance && glance-registry --config-file=/etc/glance/glance-registry.conf" glance
