# debian/glance.glance-api.upstart
description "Glance API server"
author "Soren Hansen <soren@linux2go.dk>"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

respawn

exec su -c "cd /var/lib/glance && glance-api --config-file=/etc/glance/glance-api.conf" glance
